# autostart/startup.desktop
[Desktop Entry]
Type=Application
Name=Dev Workflow Launcher
Comment=Ask which workflow to open at login
Exec=bash -c "$HOME/startup.sh"
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=8
